# kaamkartApi/kaamkart-api.service
[Unit]
Description=KaamKart API Service
Documentation=https://github.com/your-username/nokariya
After=network.target mysql.service
Requires=mysql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/kaamkart-api

# Java options
# Optimized for AWS Free Tier (t2.micro with 1GB RAM)
# Default: -Xms512m -Xmx1024m (for 2GB+ RAM)
# Free Tier: -Xms256m -Xmx384m (for 1GB RAM)
Environment="JAVA_OPTS=-Xms256m -Xmx384m -XX:+UseG1GC -XX:MaxGCPauseMillis=200"

# Spring Boot options
Environment="SPRING_PROFILES_ACTIVE=prod"

# Database configuration
Environment="DB_URL=jdbc:mysql://localhost:3306/kaamkart?useSSL=true&requireSSL=true&serverTimezone=UTC"
Environment="DB_USERNAME=kaamkart_user"
Environment="DB_PASSWORD=CHANGE_THIS_PASSWORD"

# JWT configuration
Environment="JWT_SECRET=CHANGE_THIS_SECRET"
Environment="JWT_EXPIRATION=604800000"

# CORS and WebSocket
Environment="CORS_ALLOWED_ORIGINS=https://kaamkart.in,https://www.kaamkart.in"
Environment="WEBSOCKET_ALLOWED_ORIGINS=https://kaamkart.in,https://www.kaamkart.in"

# Server port
Environment="SERVER_PORT=8585"

# Logging
Environment="LOG_FILE_PATH=/var/log/kaamkart-api.log"

# Executable
ExecStart=/usr/bin/java $JAVA_OPTS -jar -Dspring.profiles.active=prod /opt/kaamkart-api/kaamkart-api-1.0.0.jar

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log /opt/kaamkart-api

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kaamkart-api

[Install]
WantedBy=multi-user.target

